# scripts/dmr-monitor.service
[Unit]
Description=DMR Cap+ Monitor Dashboard
After=network.target

[Service]
# Adjust WorkingDirectory/ExecStart paths to your checkout location.
WorkingDirectory=/home/shahar/DMR-DECREP-SHAHAR
# Direct SDR control (no SDRconnect / virtual cable) — set --frequency to
# your Cap+ channel. Prereqs: bash scripts/setup-sdrplay.sh
# Port 8081, NOT 8080: this repo is the reference/engine source for the
# DMR ⨯ DMR-DECREP merge (see CHANGELOG "Phase 0"/"Phase 2") — DMR's own
# dmr-web.service is the always-on production app on 8080 on the same Pi.
ExecStart=/home/shahar/DMR-DECREP-SHAHAR/.venv/bin/python -m backend.cli \
    --live \
    --rf-backend soapy \
    --frequency 168.5M \
    --calls-dir /tmp/dmr_calls \
    --serve \
    --port 8081 \
    --retention-days 7
# Legacy audio-cable chain (SDRconnect + dmr_capture sink) — swap the two
# RF lines above for:
#     --input pulse:dmr_capture.monitor \
Restart=on-failure
RestartSec=10s
# Prevent a crash-loop from hammering the system: allow at most 5 restarts
# in a 2-minute window, then give up and let the operator investigate.
StartLimitBurst=5
StartLimitIntervalSec=120s
Environment=PYTHONUNBUFFERED=1
# Cap memory so a leak can't OOM the whole Pi.
MemoryMax=384M
# Log to the journal so journalctl -u dmr-monitor works.
StandardOutput=journal
StandardError=journal
# NOTE: /tmp/dmr_calls is cleaned on reboot (tmpfs / systemd-tmpfiles).
# Change --calls-dir to a persistent path such as /var/lib/dmr-monitor/calls
# if recording loss across reboots is unacceptable for your deployment.

[Install]
WantedBy=default.target
